# supervisor.conf
[program:berlin-chatbot]
command = uwsgi /usr/src/app/uwsgi.ini
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true